# ops/com.trueknot.macrodaily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.trueknot.macrodaily</string>
  <!-- --reports-dir is repeatable; newest-written run wins per ticker on a same
       trade_date (corrected reruns beat same-date originals). These dated bases
       are cadence-specific: when a NEW cadence lands in a new dir, add it here
       (or consolidate run dirs into one canonical base). --manifest-scope limits
       the plan to published (manifest) tickers so every row links. -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/trueknot/tradingagents/.venv/bin/tradingmacro</string>
    <string>--manifest-scope</string>
    <string>--reports-dir</string><string>/Users/trueknot/research-staging-2026-06-02</string>
    <string>--reports-dir</string><string>/Users/trueknot/tkruns-rerun-2026-06-02</string>
    <string>--reports-dir</string><string>/Users/trueknot/tkruns-fix-2026-06-02</string>
    <string>--reports-dir</string><string>/Users/trueknot/tkresearch</string>
    <string>--reports-dir</string><string>/Users/trueknot/.openclaw/data/research</string>
    <string>--sheet-id</string>
    <string>REPLACE_WITH_TRADING_PLAN_SHEET_ID</string>
    <string>--manifest</string>
    <string>/Users/trueknot/gsheet-tool/pdf_ids.tsv</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <!-- PATH must include /opt/homebrew/bin so the gog CLI resolves under launchd. -->
    <key>PATH</key><string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>FRED_API_KEY</key><string>REPLACE_WITH_FRED_KEY</string>
    <!-- gog Sheets write: account + keyring password (gog reads both natively). -->
    <key>GOG_ACCOUNT</key><string>trueknotsg@gmail.com</string>
    <key>GOG_KEYRING_PASSWORD</key><string>REPLACE_WITH_GOG_KEYRING_PASSWORD</string>
    <key>MACRO_CACHE_DIR</key><string>/Users/trueknot/.cache/tradingagents-macro</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>5</integer><key>Minute</key><integer>10</integer></dict>
  <key>StandardOutPath</key><string>/Users/trueknot/.macrodaily.log</string>
  <key>StandardErrorPath</key><string>/Users/trueknot/.macrodaily.err</string>
</dict>
</plist>
